# doorsd.service
[Unit]
Description=Workantile Doors RFID Daemon

[Service]
ExecStart=/usr/bin/doorsd
Restart=always

[Install]
WantedBy=multi-user.target
